# config/sys-setup.service
[Unit]
Description=Run system setup script

[Service]
ExecStart=/setup.sh

[Install]
WantedBy=multi-user.target